# crontab-entries.cron
0 2,8,14,20 * * * /etc/init.d/minecraft backup
30 * * * * /home/minecraft/cronrestart.sh
* * * * * /home/minecraft/logandgraphmemory.sh $(whoami) /home/minecraft/memorylogs memorylog.dat
